# roles/gandalf/templates/gandalf.service.j2
[Unit]
Description=Gandalf

[Service]
ExecStart=/usr/bin/gandalf-webserver
User={{ runner.name }}

[Install]
WantedBy=multi-user.target
