# roles/tsurud/templates/tsurud.service.j2
[Unit]
Description=Docker-based Platform as a Service (PaaS)

[Service]
ExecStart=/usr/bin/tsurud api
User={{ runner.name }}

[Install]
WantedBy=multi-user.target
