# templates/git-daemon.service.j2
[Unit]
Description=Git daemon

[Service]
ExecStart=/usr/bin/git daemon --base-path=/var/lib/gandalf/repositories --export-all
User=gandalf

[Install]
WantedBy=multi-user.target
