# roles/git-daemon/templates/git-daemon.service.j2
[Unit]
Description=Git daemon

[Service]
ExecStart=/usr/bin/git daemon --base-path={{ data_dir }} --export-all
User={{ runner.name }}

[Install]
WantedBy=multi-user.target
